# agent_jarvis/supervisor.d/demo_service1.ini
[program:demo_service1]
command=bash /root/jarvis/supervisor.d/demo_service1.sh
priority=1
autostart=true
autorestart=true
stdout_logfile=/root/jarvis/supervisor_log/demo_service1_stdout.log
stdout_logfile_maxbytes=4MB
stderr_logfile=/root/jarvis/supervisor_log/demo_service1_stderr.log
stderr_logfile_maxbytes=4MB