# agent_jarvis/supervisor.d/demo_service2.ini
[program:demo_service2]
command=bash /root/jarvis/supervisor.d/demo_service2.sh
priority=1
autostart=true
autorestart=true
stdout_logfile=/root/jarvis/supervisor_log/demo_service2_stdout.log
stdout_logfile_maxbytes=4MB
stderr_logfile=/root/jarvis/supervisor_log/demo_service2_stderr.log
stderr_logfile_maxbytes=4MB