# agent_jarvis/supervisor.d/jarvis.ini
[program:jarvis]
command=bash /root/jarvis/supervisor.d/jarvis.sh
priority=1
autostart=true
autorestart=true
stdout_logfile=/root/jarvis/supervisor_log/jarvis_stdout.log
stdout_logfile_maxbytes=4MB
stderr_logfile=/root/jarvis/supervisor_log/jarvis_stderr.log
stderr_logfile_maxbytes=4MB